# BTLE-Scanner/btle-scanner.service
[Unit]
Description=BTLE Scanner - HTTP Bridge
After=network-online.target bluetooth.service
Wants=network-online.target
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/local/bin/btle_scanner
Restart=on-failure
RestartSec=5
# Runs without root thanks to setcap (see Makefile "install")
User=nobody
Group=nogroup
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
